# navigator/stages/stage5_visibility_detection.json
{
  "name": "Stage 5 — Visibility, Detection & Operations",
  "versions": {"attack": "13", "navigator": "4.9", "layer": "4.5"},
  "domain": "enterprise-attack",
  "description": "ATT&CK technique coverage for Pipeline Stage 5: Visibility, Detection & Operations. Vendors: CrowdStrike Falcon, SentinelOne, Microsoft Defender XDR, Splunk ES, Microsoft Sentinel, IBM QRadar, Darktrace, ExtraHop. NIST controls: AU-2, AU-6, SI-3, SI-4, IR-4.",
  "filters": {"platforms": ["Windows", "Linux", "macOS", "Cloud", "Network"]},
  "sorting": 0,
  "layout": {"layout": "side", "aggregateFunction": "average", "showID": true, "showName": true, "showAggregateScores": false, "countUnscored": false, "expandedSubtechniques": "none"},
  "hideDisabled": false,
  "techniques": [
    {"techniqueID": "T1059", "score": 1, "comment": "Stage 5 EDR+SIEM — CrowdStrike/SentinelOne prevent; Splunk/Sentinel detect via process telemetry"},
    {"techniqueID": "T1059.001", "score": 1, "comment": "Stage 5 — PowerShell execution detection"},
    {"techniqueID": "T1059.003", "score": 1, "comment": "Stage 5 — Windows Command Shell detection"},
    {"techniqueID": "T1055", "score": 1, "comment": "Stage 5 EDR — Process injection detection via kernel telemetry"},
    {"techniqueID": "T1486", "score": 1, "comment": "Stage 5 EDR+SIEM — Ransomware behavioral detection and SOAR response"},
    {"techniqueID": "T1046", "score": 1, "comment": "Stage 5 EDR+NDR — Network scan detection from endpoint and network perspectives"},
    {"techniqueID": "T1078", "score": 1, "comment": "Stage 5 SIEM — Anomalous authentication patterns detected across identity telemetry"},
    {"techniqueID": "T1098", "score": 1, "comment": "Stage 5 SIEM — Account manipulation events correlated across directory and endpoint logs"},
    {"techniqueID": "T1190", "score": 1, "comment": "Stage 5 EDR — Post-exploitation activity detected after public-facing application exploit"},
    {"techniqueID": "T1547", "score": 1, "comment": "Stage 5 EDR — Persistence mechanism detection"},
    {"techniqueID": "T1547.001", "score": 1, "comment": "Stage 5 EDR — Registry Run Keys / Startup Folder detection"},
    {"techniqueID": "T1543", "score": 1, "comment": "Stage 5 EDR — Malicious service creation detection"},
    {"techniqueID": "T1053", "score": 1, "comment": "Stage 5 EDR+SIEM — Scheduled task creation monitoring"},
    {"techniqueID": "T1053.005", "score": 1, "comment": "Stage 5 — Scheduled Task creation by unusual parents"},
    {"techniqueID": "T1562", "score": 1, "comment": "Stage 5 EDR+SIEM — Defense impairment detection (AV disable, log clearing)"},
    {"techniqueID": "T1490", "score": 1, "comment": "Stage 5 EDR — Shadow copy deletion detection"},
    {"techniqueID": "T1489", "score": 1, "comment": "Stage 5 EDR — Service stop detection"},
    {"techniqueID": "T1071", "score": 1, "comment": "Stage 5 NDR — C2 traffic detection via Darktrace/ExtraHop behavioral analytics"},
    {"techniqueID": "T1071.001", "score": 1, "comment": "Stage 5 NDR — Unusual HTTP/S beaconing detection"},
    {"techniqueID": "T1040", "score": 1, "comment": "Stage 5 NDR — Network sniffing detection"},
    {"techniqueID": "T1048", "score": 1, "comment": "Stage 5 NDR+SIEM — Exfiltration over alternative protocol detection"},
    {"techniqueID": "T1110", "score": 1, "comment": "Stage 5 SIEM — Brute force pattern detection across authentication logs"},
    {"techniqueID": "T1136", "score": 1, "comment": "Stage 5 SIEM — Unauthorized account creation detection"},
    {"techniqueID": "T1566", "score": 1, "comment": "Stage 5 EDR — Malicious payload detection after email delivery to endpoint"},
    {"techniqueID": "T1204", "score": 1, "comment": "Stage 5 EDR — User execution of malicious content detected"}
  ],
  "gradient": {"colors": ["#ffffff00", "#ef4444ff"], "minValue": 0, "maxValue": 1},
  "legendItems": [{"label": "Covered by Stage 5 — Visibility, Detection & Operations", "color": "#ef4444ff"}],
  "metadata": [
    {"name": "Pipeline Stage", "value": "Stage 5 — Visibility, Detection & Operations"},
    {"name": "NIST Controls", "value": "AU-2, AU-6, SI-3, SI-4, IR-4"},
    {"name": "Key Vendors", "value": "CrowdStrike, SentinelOne, Splunk ES, Microsoft Sentinel, Darktrace, ExtraHop"}
  ],
  "showTacticRowBackground": false, "tacticRowBackground": "#dddddd",
  "selectTechniquesAcrossTactics": false, "selectSubtechniquesWithParent": false
}
